# stacks/litespeed-swoole/etc/service/swoole/run
#!/bin/sh

exec 2>&1
source /env

if [ ${APP} == "symfony" ]; then
  SCRIPT="/usr/local/lsws/symfony/vendor/bin/swoole --env=prod --no-debug"
else
  SCRIPT="/usr/local/lsws/Example/html/index.php"
fi

exec $SCRIPT
