# stacks/nginx-swoole-uds/etc/service/swoole/run
#!/bin/sh

exec 2>&1
source /env

if [ ${APP} == "symfony" ]; then
  SCRIPT="/var/www/symfony/vendor/bin/swoole --env=prod --no-debug"
else
  SCRIPT="/var/www/html/index.php"
fi

exec $SCRIPT
