# stacks/nginx-swoole-uds/etc/service/swoole/run
#!/bin/sh

exec 2>&1
source /env

exec php /var/www/symfony/web/swoole.php
