# stacks/nginx-swoole/etc/service/swoole/run
#!/bin/sh

exec 2>&1
source /env

exec php /var/www/app/index.php
